# src/rocks-lifemapper/lucene.cron.in
# Lifemapper lucene cron job to be installed in /etc/cron.daily

0  */4  *  *  *  @LMHOME@/rocks/bin/buildLmLuceneIndex
