# src/rocks-lifemapper/lucene.cron.in
# Lifemapper lucene cron job

0  */4  *  *  *  @PKGROOT@/rocks/bin/buildLmLuceneIndex
