# src/rocks-lifemapper/lucene.cron.in
# Lifemapper lucene cron job

0  */4  *  *  *  @LMHOME@/LmWebServer/lucene/buildLmLuceneIndex.cron
